# config/autostart/xbanish.desktop
[Desktop Entry]
Type=Application
Name=xbanish
Comment=Hide mouse cursor when typing
Exec=xbanish
NoDisplay=true
X-XFCE-Autostart-Override=true
